# components/dropbear/files/etc/uci-defaults/40-dropbear-fixes
#!/bin/sh
chmod 600 /etc/dropbear/authorized_keys
exit 0
