# components/dropbear/files/etc/uci-defaults/40-dropbear-fixes
#!/bin/sh

chmod 600 /etc/dropbear/authorized_keys

rm /root/id_rsa
rm /root/id_rsa.pub

exit 0
