# components/tor/files/etc/hotplug.d/iface/40-tor
#!/bin/sh

if [ "$INTERFACE" != "lo" ] && [ "$ACTION" = "ifup" ]
then
        /etc/init.d/tor restart
fi
