# components/jirafeau/files/etc/uci-defaults/60-jirafeau-fixes
#!/bin/sh

chmod 777 /tmp_php_*
chmod 766 -R /www/storage-jirafeau@usb/
chmod 755 -R /www/WEBtor/

exit 0
